# meta-microservicebus-node/recipes-microservicebus-node/microservicebus-node-initscript/files/microservicebus-node-initscript.service
[Unit]
Description=microServiceBus init plattform
DefaultDependencies=no
After=local-fs.target
Before=rauc.service
Before=network-pre.target
Before=microservicebus-node.service

[Service]
Type=oneshot
ExecStart=/usr/bin/microservicebus-node-initscript.sh
Restart=on-failure
RestartSec=5s
TimeoutSec=300

[Install]
WantedBy=multi-user.target
